# templates/svmp-init-rh.erb
#!/bin/bash
#
# <%= @service_name %>      Secure Virtual Mobile Platform
#
# chkconfig: 2345 80 30
# description: Startup script for the Secure Virtual Mobile Platform \
#              component <%= @daemon_name %>
# processname: <%= @daemon_name %>
# pidfile: /var/run/<%= @service_name %>.pid
# config: <%= @conf_dir %>/<%= @conf_file %>

### BEGIN INIT INFO
# Provides: <%= @service_name %>
# Required-Start: $local_fs $network $syslog
# Required-Stop: $local_fs $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: <%= @service_name %>
# Description:       Startup script for the Secure Virtual Mobile Platform
#                    component <%= @daemon_name %>
### END INIT INFO

# Author: David Keppler <dkeppler@mitre.org>

. /etc/rc.d/init.d/functions

RETVAL=0

PATH=/sbin:/usr/sbin:/bin:/usr/bin
DESC="<%= @service_name %>"
NAME=<%= @service_name %>
USER=<%= @user %>
DAEMON=/usr/bin/<%= @daemon_name %>
DAEMON_ENV="NODE_ENV=production config=<%= @conf_dir %>/<%= @conf_file %>"
FOREVER_OPTS="-a -o <%= @log_dir %>/<%= @service_name %>.out -e <%= @log_dir %>/<%= @service_name %>.err -l <%= @log_dir %>/forever.log"
DAEMON_ARGS=""
PIDFILE=/var/run/$NAME.pid
SCRIPTNAME=/etc/init.d/$NAME

do_start() {
	echo -n "Starting $NAME services: "
	su -c "$DAEMON_ENV forever start $FOREVER_OPTS $DAEMON" $USER
	touch /var/lock/subsys/<%= @service_name %>
}

do_stop() {
	echo -n "Shutting down $NAME services: "
	su -c "forever stop $DAEMON" $USER
	rm -f /var/lock/subsys/<%= @service_name %>
}

case "$1" in
	start)
		do_start
	;;
	stop)
		do_stop
	;;
	status)
	;;
	reload|restart)
		do_stop
		do_start
	;;
	*)
		echo "Usage: <servicename> {start|stop|status|reload|restart}"
		exit 1
	;;
esac
exit $RETVAL
